# deploy/trackbox.service
[Unit]
Description=TrackBox tracking API
After=network.target

[Service]
Type=simple
User=trackbox
Group=trackbox
WorkingDirectory=/opt/trackbox
# HOME must be set so Playwright can find the browser it downloaded to
# ~/.cache/ms-playwright (installed as the same user).
Environment=HOME=/opt/trackbox
# Internal port the app listens on (Nginx proxies to it). Override by setting
# APP_PORT in /opt/trackbox/.env. Default 8080 if not set.
Environment=APP_PORT=8080
EnvironmentFile=-/opt/trackbox/.env
ExecStart=/opt/trackbox/.venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 --port ${APP_PORT} --workers 1 --proxy-headers --forwarded-allow-ips=*
Restart=always
RestartSec=3
# headless Chromium can be memory hungry; give it room / restart if it leaks
MemoryMax=1500M

[Install]
WantedBy=multi-user.target
